# vm/files/startup/start-tty1.conf
start on stopped rc RUNLEVEL=[35]
stop on runlevel [S0126]

respawn
console owner
exec /usr/bin/openvt -c 1 -l -w -f -- /bin/bash /opt/startup/startup_script start

